# config/etc/services.d/ossec-logs/finish
#!/bin/bash

/var/ossec/bin/wazuh-control stop

/run/s6/basedir/bin/halt